# src/main/resources/spark-worker-supervisor.conf
[program:spark-worker]
command=sh /opt/apache-spark/sbin/start-spark-worker.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/stdout.log
autostart=true
autorestart=true